# lib/systemd/system/live-mode-apparmor.service
[Unit]
Description=Fix AppArmor for live mode

[Service]
Type=oneshot
ExecStart=/usr/lib/live-mode-apparmor

[Install]
WantedBy=multi-user.target
