# lib/systemd/system/apparmor.service.d/30_live_mode.conf
## Copyright (C) 2019 - 2019 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

## https://forums.whonix.org/t/live-mode-breaks-apparmor/7559
[Unit]
ConditionPathExists=
